# deploy/sparring.service
# systemd 单元（参考 · 授权后 cp 到 /etc/systemd/system/）
[Unit]
Description=Sparring (对练场) API
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/sparring
EnvironmentFile=/opt/sparring/.env
ExecStart=/opt/sparring/deploy/run-prod.sh
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
